# supervisord.conf
# [unix_http_server]
# file=/run/supervisor.sock
# chmod=0700

[supervisord]
pidfile=/run/supervisord.pid
nodaemon=true
loglevel=info

[supervisorctl]
# serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:sm_slider]
process_name=sm_slider
command=python3 -Bu app.py
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
